# deploy/localloop-backend.service
[Unit]
Description=LocalLoop Backend API
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/code/local-loop-io/localloop-backend
EnvironmentFile=/root/code/local-loop-io/localloop-backend/.env
ExecStart=/root/.bun/bin/bun /root/code/local-loop-io/localloop-backend/src/index.ts
Restart=on-failure
RestartSec=5

# Use a non-root user that owns the repo and data directory.
User=localloop
Group=localloop

[Install]
WantedBy=multi-user.target
