# deploy/localloop-backend.service
[Unit]
Description=LocalLoop Backend API
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/code/local-loop-io/localloop-backend
EnvironmentFile=/root/code/local-loop-io/localloop-backend/.env
ExecStart=/root/.volta/bin/node /root/code/local-loop-io/localloop-backend/src/server.js
Restart=on-failure
RestartSec=5

# Update to a non-root user that owns the repo and data directory.
User=root

[Install]
WantedBy=multi-user.target
