# deploy/localloop-backend.service
[Unit]
Description=LocalLoop Backend API
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
# Paths should be customized to match your deployment.
# Example: /opt/localloop-backend or /home/localloop/localloop-backend
WorkingDirectory=/opt/localloop-backend
EnvironmentFile=/opt/localloop-backend/.env
ExecStart=/usr/local/bin/bun run /opt/localloop-backend/src/index.ts
Restart=on-failure
RestartSec=5

# Security: Run as dedicated non-root user
# Create user with: sudo useradd -r -s /usr/sbin/nologin -d /opt/localloop-backend localloop
User=localloop
Group=localloop

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/opt/localloop-backend/data

[Install]
WantedBy=multi-user.target
